# templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
Type=simple
User={{ node_exporter__service_user }}
Group={{ node_exporter__service_group }}
Nice=-5
ExecStart={{ node_exporter__service_executable }} \
{% for collector in node_exporter__enabled_collectors -%}
{%   if not collector is mapping %}
    --collector.{{ collector }} \
{%   else -%}
{%     set name, options = (collector.items()|list)[0] -%}
{%     for k,v in options|dictsort %}
    --collector.{{ name }}.{{ k }}={{ v }} \
{%     endfor -%}
{%   endif -%}
{% endfor -%}
{% for collector in node_exporter__disabled_collectors %}
    --no-collector.{{ collector }} \
{% endfor %}
    --web.telemetry-path="{{ node_exporter__web_telemetry_path }}" \
    --web.listen-address="{{ node_exporter__web_listen_address }}"

SyslogIdentifier={{ node_exporter__service_name }}
Restart=always

[Install]
WantedBy=multi-user.target
