# templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
Type=simple
User=node-exp
Group=node-exp
ExecStart=/opt/node_exporter --web.listen-address {{ node_exporter_web_listen_address }} --collectors.enabled {{ node_exporter_enabled_collectors | join(',') }}

SyslogIdentifier=node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
